# dist/systemd/fan-control.service
[Unit]
Description = fan controller

[Service]
ExecStart = /usr/local/sbin/fan-control.sh
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target

