# dist/systemd/fan-control.service
[Unit]
Description = fan controller

[Service]
ExecStart = /usr/local/sbin/fan-control -config /usr/local/etc/fan-control.yaml
Restart = always
Type = simple

[Install]
WantedBy = multi-user.target

